# backend/labs/vulhub-net/node-db/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mariadb]
command=bash -c "mkdir -p /run/mysqld && chown mysql:mysql /run/mysqld && exec /usr/sbin/mariadbd --user=mysql --datadir=/var/lib/mysql"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf --daemonize no
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
