# backend/labs/dvwa/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

[program:mariadb]
# Recreate /run/mysqld each time (tmpfs) so the unix socket can bind.
command=bash -c "mkdir -p /run/mysqld && chown mysql:mysql /run/mysqld && exec /usr/sbin/mariadbd --user=mysql --datadir=/var/lib/mysql"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:apache]
command=/usr/local/bin/apache2-foreground
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
